# stage2/03-install-packages/files/camera360.service
[Unit]
Description=Camera360 startup script
After=network.target
After=systemd-user-sessions.service
After=network-online.target

[Service]
ExecStart=/home/pi/camera360pi
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=camera360
Restart=always
RestartSec=60


[Install]
WantedBy=multi-user.target